# git-auto-commit-postfix.service
[Unit]
Description=Auto commit Postfix configuration changes
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/git-auto-commit-postfix.sh
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
